# couchpotato.ini
[program:CouchPotatoServer]
command=/CouchPotatoServer-master/CouchPotato.py --data_dir=/config/couchpotato
