# sabnzbd.ini
[program:sabnzbd]
command=/SABnzbd-0.7.20/SABnzbd.py -f /config/sabnzbd/sabnzbd.ini
